# deploy/systemd/it-digest-watch.timer
[Unit]
Description=Run release watchers every hour

[Timer]
OnCalendar=hourly
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
